# fake-hwclock-tick.service
[Unit]
Description=Save hwclock by timer

[Service]
Type=oneshot
ExecStart=/bin/fake-hwclock tick > /dev/null

